# scripts/bayyti.service
[Unit]
Description=BAYYTI Smart Irrigation System - API Server
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/smart-irrigation/backend
Environment="PYTHONUNBUFFERED=1"
Environment="ENABLE_GPIO=true"
ExecStart=/usr/bin/python3 /home/pi/smart-irrigation/backend/api_server.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Restart conditions
StartLimitIntervalSec=60
StartLimitBurst=3

[Install]
WantedBy=multi-user.target
